# fcos-zfs-autoupdate.service
[Unit]
Description=Automatic system image updates service

[Service]
ExecStart=rpm-ostree upgrade --reboot
