# fcos-zfs-autoupdate.service
[Unit]
Description=Automatic system image updates service

[Service]
# ExecCondition=nm-online -t 10
ExecStart=rpm-ostree upgrade --reboot
